# templates/kubernetes-goat-portforward.service.j2
[Unit]
Description=Kubernetes Goat Port Forward Service
After=docker.service
Wants=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
User={{ ansible_user }}
WorkingDirectory={{ kubernetes_goat_path }}
Environment=KUBECONFIG=/home/{{ ansible_user }}/.kube/config
Environment=PATH=/usr/local/bin:/usr/bin:/bin
ExecStartPre=/bin/bash -c 'export KUBECONFIG=/home/{{ ansible_user }}/.kube/config && until /usr/local/bin/kubectl get pods -l app=kubernetes-goat-home -o jsonpath="{.items[0].status.phase}" | grep -q Running; do sleep 5; done'
ExecStart=/bin/bash {{ kubernetes_goat_path }}/access-kubernetes-goat-network.sh
ExecStop=/bin/bash -c 'if [ -f /var/run/kubernetes-goat-pids.txt ]; then while read pid; do kill $pid 2>/dev/null || true; done < /var/run/kubernetes-goat-pids.txt; rm -f /var/run/kubernetes-goat-pids.txt; fi'
StandardOutput=append:/var/log/kubernetes-goat-portforward.log
StandardError=append:/var/log/kubernetes-goat-portforward.log

[Install]
WantedBy=multi-user.target